# nsncd.service
[Unit]
Description=name-service non-caching daemon

[Service]
ExecStart=/usr/lib/nsncd
Restart=always
