# nsncd.socket
[Socket]
ListenStream=/var/run/nscd/socket

[Install]
WantedBy=sockets.target
