# nsncd.socket
[Socket]
ListenStream=/var/run/nscd/socket
